# docker/display/fakeroot/etc/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0

[program:xvnc]
priority=0
command=Xvnc -listen tcp -listen local -listen unix -geometry "%(ENV_DISPLAY_WIDTH)s"x"%(ENV_DISPLAY_HEIGHT)s" -depth 24 -securitytypes=none -rfbport="%(ENV_VNC_PORT)s" -rfbunixpath=/tmp/vnc.sock -rfbunixmode=0660
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:novnc]
priority=0
command=/opt/noVNC/utils/novnc_proxy --vnc 0.0.0.0:"%(ENV_VNC_PORT)s" --listen "%(ENV_WEB_PORT)s"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:openbox]
priority=0
command=openbox
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:configure]
priority=1
command=sh -c "/etc/configure.sh && tail -f /dev/null"
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# [program:ffmpeg]
# priority=3
# command=ffmpeg -y -f x11grab -video_size "%(ENV_DISPLAY_WIDTH)s"x"%(ENV_DISPLAY_HEIGHT)s" -r 25 -i "%(ENV_DISPLAY)s" -pix_fmt yuv420p /videos/output.mp4
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true
